# status.conf
[program:status]
command=/usr/bin/node /var/www/status/index.js
autostart=true
autorestart=true
stopasgroup=true
stderr_logfile=/var/log/status.err.log
stdout_logfile=/var/log/status.out.log
